# overlay/system/usr/lib/udev/rules.d/70-android-brightness.rules
ACTION=="add", RUN+="/bin/sh -c 'chown root:root /sys/class/backlight/panel0-backlight/*'"