# etc/init/double_tap_to_wake.conf
description "Enable double tap to wake"

start on started lxc-android-config
task

script
    echo 1 > /proc/touchpanel/double_tap_enable
end script
