# overlay/system/etc/systemd/system/repowerd.service.d/repowerd.conf
[Unit]
After=lxc-android-config.service dbus.socket sensorfwd.service

