# overlay/system/etc/systemd/system/hfd-service.service.d/hfd_service_wait_vibrator_init.conf
[Service]
ExecStartPre=/bin/bash /etc/systemd/system/hfd-service.service.d/wait_for_vibrator_init.sh

